# scripts/launchd/com.claude-web.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Stable backend ("prod"). Listens on :3030, writes to ~/.claude-web.
  Runs `pnpm start:backend` (tsx without watch) so dev edits in src/ do NOT
  auto-reload this instance. To promote new code:
      launchctl kickstart -k gui/$(id -u)/com.claude-web.backend
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude-web.backend</string>

  <!--
    2026-05-05 lesson: previous form `bash -lc "exec caffeinate -is pnpm ..."`
    had two bugs —
      (1) -lc reads login shell init which rewrites PATH (Homebrew node 25 wins
          over the nvm v24 we want), making backend run on whichever Node was
          last installed; native bindings then go NODE_MODULE_VERSION mismatch
      (2) when launchd exec'd this chain via bash, pnpm hung indefinitely under
          caffeinate (no tsx child ever spawned, log empty)
    Fix: keep /bin/bash as the launchd program (preserves macOS TCC permission
    for ~/Desktop access — switching to /pnpm would require fresh Full Disk
    Access grant), but drop -l and caffeinate.
    Tradeoff: lose caffeinate idle-sleep prevention; acceptable for personal
    Mac that's plugged in, plus KeepAlive=true respawns if Mac wakes.
  -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-c</string>
    <string>exec /Users/yongqian/.nvm/versions/node/v24.12.0/bin/pnpm start:backend</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/yongqian/Desktop/claude-web-prod</string>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- ~/.local/bin must be on PATH so scheduler-spawned `claude` CLI is
         findable. Was previously injected by login-shell init when plist used
         `bash -lc`; now -l is dropped, the path must be explicit. -->
    <key>PATH</key>
    <string>/Users/yongqian/.local/bin:/Users/yongqian/.nvm/versions/node/v24.12.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/yongqian</string>
    <key>PORT</key>
    <string>3030</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>StandardOutPath</key>
  <string>/Users/yongqian/Library/Logs/claude-web-backend.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/yongqian/Library/Logs/claude-web-backend.stderr.log</string>
</dict>
</plist>
