# scripts/launchd/com.claude-web.backend.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Stable backend ("prod"). Listens on :3030, writes to ~/.claude-web.
  Runs `pnpm start:backend` (tsx without watch) so dev edits in src/ do NOT
  auto-reload this instance. To promote new code:
      launchctl kickstart -k gui/$(id -u)/com.claude-web.backend
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.claude-web.backend</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>exec /usr/bin/caffeinate -is /Users/yongqian/.nvm/versions/node/v24.12.0/bin/pnpm start:backend</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/yongqian/Desktop/claude-web</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/yongqian/.nvm/versions/node/v24.12.0/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    <key>HOME</key>
    <string>/Users/yongqian</string>
    <key>PORT</key>
    <string>3030</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>ThrottleInterval</key>
  <integer>5</integer>

  <key>StandardOutPath</key>
  <string>/Users/yongqian/Library/Logs/claude-web-backend.stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/yongqian/Library/Logs/claude-web-backend.stderr.log</string>
</dict>
</plist>
